# docker-kubernetes/logging-system/log4j2 systemd journal log fluent bit es/poc/app-file/log4j2-app.service
[Unit]
Description=Log4j2 file appender PoC
After=systemd-journald.service
Wants=systemd-journald.service

[Service]
Type=simple
ExecStartPre=/bin/sh -c "test -s /etc/machine-id || systemd-machine-id-setup"
ExecStartPre=/bin/sh -c "mkdir -p /var/log/log4j2-app"
ExecStart=/usr/bin/java -cp "/opt/log4j2-app/classes:/opt/log4j2-app/lib/*" -Dlog4j.configurationFile=/opt/log4j2-app/conf/log4j2.xml Main
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
